# rules/99-dongle-auto-restart.rules
# 99-dongle-auto-restart.rules — Auto-restart dongle + fix permissions after USB hotplug
# Fires when any ttyUSB device is added:
#   1. Sets permissions so Asterisk can open the device
#   2. Triggers the auto-restart script in background (debounce + 5s delay + per-dongle restart)
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyUSB*", MODE="0666", GROUP="dialout", RUN+="/bin/bash -c '/opt/issabel-dashboard/scripts/auto-restart-dongle.sh %k &'"
